# setup/jenkins/autossh.conf
# autossh startup Script

description "autossh daemon startup to tunnel to jenkins"

start on net-device-up IFACE=eth0
stop on runlevel [01S6]

respawn
respawn limit 0 10 # respawn unlimited times every 20 seconds

script
export AUTOSSH_PIDFILE=/var/run/autossh.pid
export AUTOSSH_POLL=60
export AUTOSSH_FIRST_POLL=30
export AUTOSSH_GATETIME=0
export AUTOSSH_DEBUG=1
autossh -M 0 -f -N -o ServerAliveInterval=30 -o ServerAliveCountMax=3 -o BatchMode=yes -o StrictHostKeyChecking=no -i /root/.ssh/id_rsa -L8080:148.251.110.18:8080 -L6001:148.251.110.18:6001 -L29418:148.251.110.21:29418 root@148.251.110.18
end script
