# setup/jenkins/jenkins_slave_start.conf
description     "Jenkins Slave"
author          "OpenContrail CI Admin (ci-admin@opencontrail.org)"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

respawn
respawn limit 0 10

env HOME=/root
umask 007

# The default of 5 seconds is too low for mysql which needs to flush buffers
kill timeout 300

script
#   wget "http://maven.jenkins-ci.org/content/repositories/releases/org/jenkins-ci/plugins/swarm-client/1.15/swarm-client-1.15-jar-with-dependencies.jar"
    SLAVE_NAME=`cat /etc/hostname  | cut -d '.' -f 1`
    java -jar /root/swarm-client-1.15-jar-with-dependencies.jar -labels juniper-tests -mode normal -master http://jenkins.opencontrail.org:8080/ -fsroot /home/jenkins -username ci-admin -password b8057c342d44883f750d93f1cc2d092f -name $SLAVE_NAME
#   exec java -jar /root/slave.jar -jnlpUrl https://jenkins.opencontrail.org/computer/$SLAVE_NAME/slave-agent.jnlp -jnlpCredentials ci-admin:b8057c342d44883f750d93f1cc2d092f -slaveLog /root/$SLAVE_NAME.log
end script

