# modules/solr/files/upstart.conf
description "Solr search service"
author      "Marton Kiss <marton.kiss@gmail.com>"

start on runlevel [2345]
stop on runlevel [016]

kill timeout 30
respawn
script
  chdir /srv/solr/example
  exec /usr/bin/java -jar start.jar -Dsolr.solr.home=/etc/solr
end script