# modules/jenkins/files/slave_scripts/etc_init/autossh.conf
# autossh startup Script

description "autossh daemon startup to tunnel to jenkins"

start on net-device-up IFACE=eth0
stop on runlevel [01S6]

respawn
respawn limit 0 20 # respawn unlimited times every 20 seconds

script
export AUTOSSH_PIDFILE=/var/run/autossh.pid
export AUTOSSH_POLL=60
export AUTOSSH_FIRST_POLL=30
export AUTOSSH_GATETIME=0
export AUTOSSH_DEBUG=1
autossh -M 0 -N -o ServerAliveInterval=30 -o ServerAliveCountMax=3 -o BatchMode=yes -o StrictHostKeyChecking=no -i /root/.ssh/id_rsa -L*:29418:review.openstack.org:29418 root@jenkins.opencontrail.org
end script
