# modules/jenkins/files/slave_scripts/etc_init/slave_start.conf
description     "Jenkins Slave"
author          "OpenContrail CI Admin (ci-admin@opencontrail.org)"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

respawn
respawn limit 0 10

env HOME=/root
umask 007
kill timeout 300

script
    # Don't register in test mode.
    HOSTNAME=`cat /etc/hostname`
    ssh root@jenkins.opencontrail.org ls /root/ci-test/\*$HOSTNAME\*-test
    if [ "$?" == "0" ]; then
        echo Skip register due to jenkins.opencontrail.org:/root/ci-test/\*$HOSTNAME\*-test
        exit 0
    fi
    ruby /usr/local/jenkins/slave_scripts/ci-infra/slave_auto_run.rb
end script
