# modules/openstack_project/templates/grafana/helga.service.erb
[Unit]
Description = Helga

[Service]
Type = simple
ExecStart = /usr/helga-env/bin/helga --settings=/etc/helga/settings.py
User = helga
Group = helga